# roles/cantaloupe/templates/cantaloupe.service.j2
[Unit]
Description=Cantaloupe IIIF server
After=syslog.target network.target

StartLimitIntervalSec=300s
StartLimitBurst=3

[Service]
SuccessExitStatus=143

User=cantaloupe
Group=cantaloupe

Type=simple

WorkingDirectory=/opt/cantaloupe-{{ version }}
ExecStart=/usr/bin/java -Dcantaloupe.config=cantaloupe.properties -Xmx{{ max_heap_space }} -jar cantaloupe-{{ version }}.jar
ExecStop=/bin/kill -15 $MAINPID

Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
